# profile.d/extra.symlink
# Last things sourced save /etc/bashrc

GIT_AUTHOR_NAME="Itay Neeman"
GIT_COMMITTER_NAME="$GIT_AUTHOR_NAME"
git config --global user.name "$GIT_AUTHOR_NAME"
GIT_AUTHOR_EMAIL="itay@neeman.net"
GIT_COMMITTER_EMAIL="$GIT_AUTHOR_EMAIL"
git config --global user.email "$GIT_AUTHOR_EMAIL"

if [[ $(uname -s) == "Darwin" ]]; then
    source ~/.profile.osx.d/extra
elif [[ $(uname -s) == "Linux" ]]; then
    source ~/.profile.linux.d/extra
fi
